# supervisor/conf.d/woami100.conf
[program:whoami100]
command=/whoami --port 100