# supervisor/conf.d/woami110.conf
[program:whoami110]
command=/whoami --port 110