# supervisor/conf.d/woami90.conf
[program:whoami90]
command=/whoami --port 90