# supervisor/conf.d/woami80.conf
[program:whoami80]
command=/whoami --port 80